# templates/postgresql.service.j2
# ansible managed
.include /lib/systemd/system/postgresql.service
[Service]
Environment=PGDATA={{ postgresql_data_directory }}
